# polarisfolio-update.timer
[Unit]
Description=PolarisFolio auto-update poller (git pull + restart if changed)

[Timer]
# Poll GitHub every 10 minutes. update.sh is a no-op when nothing changed, so
# this is the primary deploy mechanism here (the webhook is blocked by
# Cloudflare Access). First run shortly after boot, then every 10 min.
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
