# polarisfolio-update.service
[Unit]
Description=PolarisFolio auto-update (git pull + restart if changed)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as root: the container's sudo is unavailable, so update.sh restarts the
# service directly when root (and chowns the tree back to the app user). Adjust
# the path if the app lives elsewhere (deploy.sh installs under the app user).
WorkingDirectory=/home/polarisfolio/polarisfolio_web
ExecStart=/home/polarisfolio/polarisfolio_web/update.sh
